# etc/systemd/system/metar-v4.service.d/logging.conf
[Service]
LogRateLimitIntervalSec=30s
LogRateLimitBurst=500
